# meta-ntb/recipes-kernel/flink/flink/flink.service
[Unit]
Description=Initialize flink FPGA

[Service]
Type=oneshot
ExecStart=/usr/bin/initflink

[Install]
WantedBy=multi-user.target