# meta-ntb/recipes-connectivity/ethernet/udhcpc-service/udhcpc.service
[Unit]
Description=Start udhcpc DHCP client

[Service]
Type=oneshot
ExecStart=/sbin/udhcpc

[Install]
WantedBy=network.target
